# conf/supervisor/conf.d/supervisor.conf
[supervisord]
nodaemon=true

[program:nginx]
command=/usr/local/nginx-streaming/sbin/nginx
stdout_logfile=/var/log/%(program_name)s.log
stderr_logfile=/var/log/%(program_name)s.log